# deploy/polymarket-dashboard.service
[Unit]
Description=Polymarket Streamlit Dashboard (read-only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/polymarket
EnvironmentFile=/home/ubuntu/polymarket/.env
Environment=PYTHONUNBUFFERED=1
Environment=STREAMLIT_BROWSER_GATHER_USAGE_STATS=false
ExecStart=/home/ubuntu/polymarket/venv/bin/streamlit run \
    /home/ubuntu/polymarket/tools/dashboard/app.py \
    --server.address 0.0.0.0 \
    --server.port 8501 \
    --server.headless true \
    --server.fileWatcherType none \
    --browser.gatherUsageStats false
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/polymarket/logs/dashboard.log
StandardError=append:/home/ubuntu/polymarket/logs/dashboard.err

[Install]
WantedBy=multi-user.target
